# rc.local
#!/bin/sh -e

/home/pi/sisbot-server/sisbot/startup.sh >> /var/log/sisyphus/startup.log  2>&1 &
exit 0
