# rc.local
#!/bin/sh -e

{
  /home/pi/sisbot-server/sisbot/startup.sh >> /var/log/sisyphus/startup.log  2>&1
} || {
  # fallback to old method
  cd /home/pi/sisbot-server/sisproxy && git reset --hard && sudo NODE_ENV=sisbot node server.js >> /var/log/sisyphus/proxy.log  2>&1 &
}
exit 0
