# rc.local
#!/bin/sh -e
cd /home/pi/sisbot-server/sisproxy && NODE_ENV=sisbot node server.js >> /var/log/sisyphus/proxy.log  2>&1
exit 0
